# admin/scripts/etc/cron.d/cloud_swg_admin_backup
#
# for development run backup every hour at HH:05 - in future let it run daily at say 01:10 am
#
5 * * * * proxy /bin/bash /opt/cloud-swg-admin/bin/backup.sh 2>&1 >/opt/cloud-swg-admin/var/log/"backup_$(date +\%Y\%m\%d_\%H\%M\%S).log"
